# sample-app-2/roles/scheduler/files/update/update_packages.service
[Unit]
Description=Service to update packages.

[Service]
Type=oneshot
ExecStart=/bin/bash /root/update_packages.sh