# sample-app-2/roles/scheduler/files/update_packages.timer
[Unit]
Description=Runs every 1 minutes test.sh

[Timer]
OnCalendar=*:0/1

[Install]
WantedBy=timers.target